# cron.daily
# From the https://wiki.archlinux.org/index.php/Reflector website

/etc/cron.daily/mirrorlist

#!/bin/bash

# Get the country thing
/usr/bin/reflector -c "Netherlands" -p http --sort rate > /etc/pacman.d/mirrorlist

# Work through the alternatives
/usr/bin/reflector -p http  --latest 20 -p https -p ftp --sort rate >> /etc/pacman.d/mirrorlist
